# consul.conf
[program:consul_server]
command=  /usr/local/bin/consul agent -server -bootstrap-expect 1 -data-dir /etc/vault/vault_files/ -bind 127.0.0.1
priority=1
autostart=true
autorestart=true
stderr_logfile=/var/log/consul.err.log
stdout_logfile=/var/log/consul.stdout.log
